# files/etc/uci-defaults/99-del-diy-feeds
#!/bin/sh

sed -i '/immortalwrt_fanx/d' /etc/opkg/distfeeds.conf
