# files/etc/uci-defaults/99-del-diy-feeds
#!/bin/sh

sed -i '/fanx/d' /etc/apk/repositories.d/distfeeds.list
