# files/etc/uci-defaults/zzz-openclash-restore-enable
#!/bin/sh

[ -x /etc/init.d/openclash ] || exit 0
[ -f /etc/config/openclash ] || exit 0

# Read pre-upgrade enable state from OpenClash backup config.
old_enable="$(uci -q -c /usr/share/openclash/backup get openclash.config.enable 2>/dev/null)"
[ "$old_enable" = "1" ] || exit 0

cfg="$(uci -q get openclash.config.config_path 2>/dev/null)"
[ -n "$cfg" ] || exit 0
[ -f "$cfg" ] || exit 0

uci -q set openclash.config.enable='1'
uci -q commit openclash

# Keep service enabled at boot.
/etc/init.d/openclash enabled >/dev/null 2>&1 || /etc/init.d/openclash enable >/dev/null 2>&1

exit 0
