# files/etc/uci-defaults/zzzz-openclash-sysupgrade-keep
#!/bin/sh

# OpenClash's default sysupgrade keep rule can preserve the whole
# /etc/openclash directory. On x86, sysupgrade stores the config archive on
# the small boot partition before reboot, so backing up Geo databases,
# caches, history, and duplicate runtime files can make sysupgrade.tgz too
# large to restore. Keep only user configuration plus the active Meta core so
# OpenClash can start immediately after upgrade without prompting to download
# the core again.

mkdir -p /lib/upgrade/keep.d

cat > /lib/upgrade/keep.d/luci-app-openclash <<'KEEP_EOF'
/etc/config/openclash
/etc/openclash/config/
/etc/openclash/custom/
/etc/openclash/overwrite/
/etc/openclash/core/clash_meta
KEEP_EOF

exit 0
