# files/etc/uci-defaults/zzzz-openclash-sysupgrade-keep
#!/bin/sh

# Keep OpenClash backups small enough for x86 sysupgrade's boot partition.

mkdir -p /lib/upgrade/keep.d

cat > /lib/upgrade/keep.d/luci-app-openclash <<'KEEP_EOF'
/etc/config/openclash
/etc/openclash/config/
/etc/openclash/custom/
/etc/openclash/overwrite/
/etc/openclash/core/clash_meta
KEEP_EOF

exit 0
