# files/etc/hotplug.d/iface/51-miniupnpd-external_ip
#!/bin/bash
sleep 60

if [ "$ACTION" = "ifup" ] && [ "$INTERFACE" = "wan" ]; then
    /usr/bin/update-upnp.sh
fi
